# docs/nixos/obscura-sidecar.service
[Unit]
Description=Obscura Stealth Browser Sidecar (50V3R31GN-M4CH1N4 Node C) — CDP on :9222
After=network.target

[Service]
Type=simple
User=maczz
WorkingDirectory=/home/maczz/50V3R31GN-M4CH1N4

# Obscura serve exposes a CDP WebSocket on --port 9222.
# Compiled with --features stealth for TLS fingerprint spoofing +
# canvas/GPU randomization + 3,520 tracker domain blocking.
ExecStart=/home/maczz/50V3R31GN-M4CH1N4/result/bin/obscura serve --port 9222

Restart=always
RestartSec=5
StandardOutput=append:/home/maczz/50V3R31GN-M4CH1N4/data/logs/obscura-sidecar.log
StandardError=append:/home/maczz/50V3R31GN-M4CH1N4/data/logs/obscura-sidecar.log

[Install]
WantedBy=multi-user.target
